# upstart/mysite.conf
description "my site (ianww.com.com)"
start on filesystem and net-device-up IFACE=eth0
stop on shutdown

exec sudo -u ian "/home/ian/my-site/upstart/upstart.sh" >>/var/log/upstart/mysite.log
